# cfpydev-imagr/i_i_aws_deploy/supervisord.conf
[program:lowr]
command: /usr/bin/python -m gunicorn imagr_site.wsgi:application
directory: /home/ubuntu
user: ubuntu
autostart: true
environment : %(host_envs)s